# mnt-shared-propagation.service
[Unit]
Description=Ensure /mnt has shared mount propagation for LXC visibility
DefaultDependencies=no
After=local-fs.target
Before=lxc.service pve-guests.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/mount --make-rshared /mnt

[Install]
WantedBy=multi-user.target
